# deploy/systemd/slskd.service
[Unit]
Description=slskd — Soulseek daemon for Muzika
Documentation=https://github.com/slskd/slskd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Same UID as muzika.service so both services share the downloads dir
# (/srv/muzika/data/music) without a permission dance.
User=muzika
Group=muzika
WorkingDirectory=/var/lib/slskd
EnvironmentFile=/etc/muzika/slskd.env
ExecStart=/opt/slskd/slskd --config /etc/muzika/slskd.yml --app /var/lib/slskd
Restart=on-failure
RestartSec=10
# cgroup v2 memory caps.
MemoryMax=400M
MemoryHigh=350M

# Sandboxing.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/var/lib/slskd /srv/muzika/data/music

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
